# Resources/launchd/ai.mnemo.stack.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.mnemo.stack</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/m3-max/Documents/6767/.build/debug/mnemoctl</string>
        <string>start</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/m3-max/Documents/6767</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/ai.mnemo.stack.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/ai.mnemo.stack.err</string>
</dict>
</plist>
